# auth_system/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server
autorestart=true

[program:django]
command=gunicorn auth_system.wsgi:application --bind 0.0.0.0:8000
autorestart=true
